# systemd/codex-pool.service
[Unit]
Description=Codex Pool
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__CODEX_POOL_WORKDIR__
ExecStart=__CODEX_POOL_BIN__ --pool-dir __CODEX_POOL_WORKDIR__/pool
Environment=CODEX_POOL_CONFIG=__CODEX_POOL_CONFIG__
Restart=on-failure
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=default.target
